# docs/examples/kaminski-service-worker.service
[Unit]
Description=Kaminski Service Solid Queue Worker
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/var/www/kaminski-service/current
EnvironmentFile=/etc/kaminski-service/env
ExecStart=/usr/bin/env bin/jobs
Restart=always

[Install]
WantedBy=multi-user.target
